# deploy/templates/palsy-backend.service.j2
[Unit]
Description=Palsy Video Training Backend
After=network.target

[Service]
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_src_path }}/backend
EnvironmentFile={{ app_root }}/backend.env
ExecStart={{ venv_path }}/bin/uvicorn app.main:app --host 0.0.0.0 --port {{ uvicorn_port }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
